# mods/etc/profile.d/firefox_wayland.sh
# Enable firefox wayland support
export MOZ_ENABLE_WAYLAND=1
