# mods/etc/profile.d/qt.sh
export QT_QPA_PLATFORM="xcb"
export QT_STYLE_OVERRIDE="kvantum"
export QT_QPA_PLATFORMTHEME="qt5ct"