# mods/etc/profile.d/wayland.sh
export QT_QPA_PLATFORM="wayland"
export QT_STYLE_OVERRIDE="kvantum"
export QT_QPA_PLATFORMTHEME="qt5ct"
export GDK_BACKEND=wayland
export MOZ_ENABLE_WAYLAND=1
export SDL_VIDEODRIVER=wayland
export CLUTTER_BACKEND=wayland
export COGL_RENDERER=egl_wayland